# agent-dist/com.lazarus.sovereign-runtime-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- The wrapper reads the bearer token from an owner-only file. Secret
     material is intentionally never embedded in this plist. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lazarus.sovereign-runtime-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>__AGENT_WRAPPER__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/agent.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/agent.log</string>
</dict>
</plist>
